# cookbooks/billy/templates/default/billy_web.sd.conf.erb
# This file was generated by Chef for '<%= node['fqdn'] %>'.
# Do NOT edit this file by hand!

[program:billy_web]
command=<%= node[:billy][:install_dir] %>/env/bin/uwsgi --yaml billy_web.uwsgi.yaml ;
process_name=%(program_name)s ;
directory=<%= node[:billy][:install_dir] %> ;
priority=999 ;
stopsignal=QUIT ;
redirect_stderr=true ;
stdout_logfile=/home/<%= node[:billy][:user] %>/logs/billy_web_out.txt ;
stderr_logfile=/home/<%= node[:billy][:user] %>/logs/billy_web_err.txt ;
