# cookbooks/billy/templates/default/web.sd.conf.erb
# This file was generated by Chef for '<%= node['fqdn'] %>'.
# Do NOT edit this file by hand!

[program:billy_web]
command=<%= node[:billy][:install_dir] %>/env/bin/uwsgi
    --home=env
    --pidfile=web.uwsgi.pid
    --socket=/tmp/billy_web.sock
    --chmod-socket=666
    --master
    --no-orphans
    --processes=8
    --max-requests=10249
    --module=wsgi
;
process_name=%(program_name)s ;
directory=<%= node[:billy][:install_dir] %> ;
priority=999 ;
stopsignal=QUIT ;
redirect_stderr=true ;
stdout_logfile=/home/<%= node[:billy][:user] %>/logs/billy_web_out.txt ;
stderr_logfile=/home/<%= node[:billy][:user] %>/logs/billy_web_err.txt ;
